# systemd/bwsshd.service
[Unit]
Description=bwsshd - keep ssh_config in sync with the Bitwarden SSH agent
After=default.target

[Service]
ExecStart=%h/.local/bin/bwsshd -watch
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
UMask=0077

[Install]
WantedBy=default.target
